# arch-desktop/services/auto-hibernate.timer
[Unit]
Description=Timer for automatic hibernation

[Timer]
OnUnitInactiveSec=30min
Unit=auto-hibernate.service

[Install]
WantedBy=timers.target
